# config/runit/tvheadend/run
#!/bin/sh
exec /usr/bin/tvheadend -C -u tvheadend -g tvheadend
